# ops/systemd/pipeline-nav-snapshot.timer
[Unit]
Description=Nightly portfolio NAV snapshot timer

[Timer]
OnCalendar=*-*-* 16:10:00
Persistent=true

[Install]
WantedBy=timers.target
